# metal/roles/loop-device/templates/rc.local.j2
#!/bin/bash

{% if loop_disk_device == '' -%}
losetup {{ loop_device }} {{ loop_file }}
{% else -%}
losetup {{ loop_device }} {{ loop_disk_device }}
{% endif %}
exit 0
