# metal/roles/loop-device/templates/rc.local.j2
#!/bin/bash

losetup {{ loop_device }} {{ loop_file }}

exit 0
